# scripts/ps3fix.timer
[Unit]
Description=Run ps3ify-sweep every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
