# scripts/ps3fix.service
[Unit]
Description=PS3 audio/container normalization for fresh imports (movie-server)
After=docker.service

[Service]
Type=oneshot
User=brennan
ExecStart=/home/brennan/movie-server/scripts/ps3ify-sweep.sh
