# scripts/ps4fix.timer
[Unit]
Description=Run ps4ify-sweep every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
